# build/osx/org.jivesoftware.openfire.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>WILDFIRE_HOME</key>
		<string>/usr/local/wildfire/</string>
	</dict>
	<key>Label</key>
	<string>org.jivesoftware.openfire</string>
	<key>OnDemand</key>
	<false/>
	<key>Program</key>
	<string>/usr/local/openfire/bin/extra/openfire-launchd-wrapper.sh</string>
	<key>RunAtLoad</key>
	<false/>
	<key>ServiceDescription</key>
	<string>Openfire</string>
	<key>StandardErrorPath</key>
	<string>/usr/local/openfire/logs/launchd.stderr</string>
	<key>StandardOutPath</key>
	<string>/usr/local/openfire/logs/launchd.stdout</string>
	<key>UserName</key>
	<string>openfire</string>
	<key>WorkingDirectory</key>
	<string>/usr/local/openfire/lib</string>
</dict>
</plist>
